# autostart.sh
#!/bin/bash
# chkconfig: 345 90 90
# description: autostart......
python /root/test.py
